# services/ossuary-web.service
[Unit]
Description=Ossuary Web Configuration Interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/ossuary/scripts/config-server.py
Restart=always
RestartSec=10
User=root
WorkingDirectory=/opt/ossuary

# Allow binding to port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target