# systemd/ossuary-portal.service
[Unit]
Description=Ossuary Portal Server
After=ossuary-api.service network-online.target
Wants=network-online.target

[Service]
Type=exec
User=root
Group=root
ExecStart=/opt/ossuary/venv/bin/python /opt/ossuary/bin/ossuary-portal
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/ossuary/src

[Install]
WantedBy=multi-user.target