# systemd/ossuary-api.service
[Unit]
Description=Ossuary API Gateway
After=ossuary-netd.service
Before=ossuary-portal.service ossuary-kiosk.service

[Service]
Type=exec
User=ossuary
Group=ossuary
ExecStart=/opt/ossuary/venv/bin/python /opt/ossuary/bin/ossuary-api
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target