# systemd/ossuary-api.service
[Unit]
Description=Ossuary API Gateway
After=ossuary-config.service ossuary-netd.service
Before=ossuary-portal.service ossuary-kiosk.service

[Service]
Type=exec
User=root
Group=root
ExecStart=/opt/ossuary/venv/bin/python /opt/ossuary/bin/ossuary-api
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/ossuary/src
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/ossuary/.Xauthority

[Install]
WantedBy=multi-user.target