# systemd/ossuary-display.service
[Unit]
Description=Ossuary Display Manager
After=ossuary-config.service graphical.target
Before=ossuary-kiosk.service
Wants=graphical.target

[Service]
Type=exec
User=root
Group=root
ExecStart=/opt/ossuary/venv/bin/python /opt/ossuary/bin/ossuary-display
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/ossuary/src

# Allow access to video devices and input
SupplementaryGroups=audio video input tty

[Install]
WantedBy=multi-user.target