# systemd/ossuary-kiosk.service
[Unit]
Description=Ossuary Kiosk Browser
After=ossuary-api.service ossuary-display.service network-online.target
Wants=network-online.target
Requisite=ossuary-display.service

[Service]
Type=exec
User=root
Group=root
ExecStart=/opt/ossuary/venv/bin/python /opt/ossuary/bin/ossuary-kiosk
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/ossuary/src

# Allow access to X display and audio
SupplementaryGroups=audio video input ossuary

[Install]
WantedBy=multi-user.target