# systemd/ossuary-wifi-monitor.service
[Unit]
Description=Ossuary WiFi Connection Monitor
After=network.target wpa_supplicant.service
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/ossuary/services/wifi_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target