# systemd/ossuary-kiosk.service
[Unit]
Description=Ossuary Kiosk Browser
After=ossuary-api.service graphical.target network-online.target
Wants=graphical.target network-online.target
Requisite=graphical.target

[Service]
Type=exec
User=root
Group=root
ExecStart=/opt/ossuary/venv/bin/python /opt/ossuary/bin/ossuary-kiosk
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/ossuary/.Xauthority
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/ossuary/src

# Allow access to X display and audio
SupplementaryGroups=audio video input ossuary

[Install]
WantedBy=graphical.target